# debian/postinst
useradd -m --system url-bot-rs --shell /usr/sbin/nologin
systemctl enable url-bot-rs.service
