# debian/prerm
systemctl disable --now url-bot-rs.service
